# monitoring/mariadb-ping.socket
[Unit]
Description=MariaDB monitoring script socket

[Socket]
ListenStream=127.0.0.1:8080
Accept=yes

[Install]
WantedBy=sockets.target
