# monitoring/mariadb-ping.service
[Unit]
Description=MariaDB monitoring script

[Service]
Type=simple
ExecStart=/usr/local/bin/mariadb-status.sh

[Install]
WantedBy=multi-user.target
